# deploy/717-scorekeeper.service
[Unit]
Description=Seven One Seven Scorekeeper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=paul
Group=paul
WorkingDirectory=/home/paul/projects/717-scorekeeper
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=7170
Environment=PUBLIC_HOST=717.novacovici.dev
Environment=DB_PATH=/home/paul/projects/717-scorekeeper/data/717.db
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5s
TimeoutStopSec=15s
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/paul/projects/717-scorekeeper/data
InaccessiblePaths=/home/paul/.cloudflared -/home/paul/.ssh -/home/paul/.gnupg
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
RemoveIPC=true
ProtectProc=invisible
ProcSubset=pid
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
